# voxel.service
[Unit]
Description=Voxel Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/voxel
ExecStart=/home/pi/voxel/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
