# services/voxel.service
[Unit]
Description=Voxel Backend (WebSocket server, hardware I/O, AI pipelines)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/voxel
ExecStart=/home/pi/.local/bin/uv run server.py
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voxel

[Install]
WantedBy=multi-user.target
