# services/voxel-web.service
[Unit]
Description=Voxel Remote UI (static web server for remote browsers)
After=voxel.service
Wants=voxel.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/voxel/app/dist
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 -m http.server 8081 --bind 0.0.0.0 -d /home/pi/voxel/app/dist

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=voxel-web

[Install]
WantedBy=multi-user.target
