# services/voxel-ui.service
[Unit]
Description=Voxel UI (WPE/Cog browser rendering React app on LCD)
After=voxel.service
Wants=voxel.service

[Service]
Type=simple
User=pi
RuntimeDirectory=voxel-ui
RuntimeDirectoryMode=0700
Environment=XDG_RUNTIME_DIR=/run/voxel-ui
Environment=COG_PLATFORM=drm
Environment=COG_PLATFORM_DRM_RENDERER=gles2
Environment=WPE_BCMRPI_TOUCH=1

# Load the pre-built React app directly from disk.
# The app connects to ws://localhost:8080 for the backend (falls back gracefully).
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/cog file:///home/pi/voxel/app/dist/index.html

Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voxel-ui

[Install]
WantedBy=multi-user.target
